# config/supervisord.conf
[program:reactor]
command=/usr/bin/node statsd/stats.js config/reactor.js
directory=/var/www/reactor/reactor/
autostart=true
autorestart=true
startsecs=2
startretries=3
user=yarek
